# docker/supervisor-nuxt.conf
[program:nuxt-server]
command = /usr/local/bin/nuxt-server.sh
startsecs = 10
autorestart = false
startretries = 1
